# modules/plugins/cron/src/test/resources/etc/cron.d/another-crontab
#another crontab file loaded from under cron.d
SHELLVAR=value

1  0  *  *  * root echo "tmp"
