# modules/plugins/cron/src/test/resources/etc/crontab
#a testing crontab file
0  0  1  1  * root echo "tmp"
